# php/7.2-fpm-alpine-nginx/config/s6/nginx/finish
#!/bin/sh
## some killall commands to stop service may be

kill `cat /var/run/nginx.pid`

## Add your cleanup tasks above this line
echo "finish: nginx"